# host_tools/sd_image/rootfs/lib/systemd/system/chameleon96-xfce-init.service
[Unit]
Description=Arrow Chameleon96 XFCE Session
After=chameleon96-fpga-init.service

[Service]
Type=simple
ExecStart=/usr/bin/startxfce4

[Install]
WantedBy=basic.target
