# host_tools/sd_image/rootfs/var/lib/opkg/info/busybox-syslog.list
/etc/syslog.conf.busybox
/lib/systemd/system/busybox-syslog.service
/etc/default/busybox-syslog
/lib/systemd/system/syslog.service
/lib/systemd/system/busybox-klogd.service
