# host_tools/sd_image/rootfs/var/lib/opkg/info/chameleon96-x11vnc-init.list
/lib/systemd/system/chameleon96-x11vnc-init.service
/usr/bin/chameleon96-x11vnc-init.sh
