# host_tools/sd_image/rootfs/var/lib/opkg/info/chameleon96-xfce-init.list
/lib/systemd/system/chameleon96-xfce-init.service
