# host_tools/sd_image/rootfs/var/lib/opkg/info/chameleon96-fpga-init.list
/usr/bin/chameleon96-fpga-init.sh
/lib/systemd/system/chameleon96-fpga-init.service
