# host_tools/sd_image/rootfs/var/lib/opkg/info/lighttpd.list
/www/var
/lib/systemd/system/lighttpd.service
/www/pages/dav/
/www/logs
/etc/lighttpd.conf
/usr/sbin/lighttpd
/usr/sbin/lighttpd-angel
/www/
/etc/lighttpd.d/
/www/pages/index.html
/www/pages/
