# host_tools/sd_image/rootfs/var/lib/opkg/info/nfs-utils-client.list
/var/lib/nfs/rmtab
/sbin/mount.nfs4
/usr/sbin/showmount
/etc/nfs-utils.conf
/var/lib/nfs/statd/sm.bak/
/usr/sbin/rpc.idmapd
/sbin/umount.nfs4
/sbin/mount.nfs
/etc/nfsmount.conf
/var/lib/nfs/v4recovery/
/usr/sbin/start-statd
/var/lib/nfs/xtab
/var/lib/nfs/
/var/lib/nfs/statd/
/var/lib/nfs/statd/sm/
/var/lib/nfs/statd/state
/usr/sbin/nfsstat
/var/lib/nfs/etab
/lib/systemd/system/nfs-statd.service
/sbin/umount.nfs
/usr/sbin/rpc.statd
/usr/sbin/sm-notify
