# host_tools/sd_image/rootfs/var/lib/opkg/info/ppp.list
/usr/sbin/pppd
/etc/chatscripts/pap
/etc/ppp/pap-secrets
/etc/ppp/ip-down.d/
/etc/ppp/ip-up.d/08setupdns
/etc/ppp/
/etc/ppp/ip-down.d/92removedns
/etc/ppp/ip-down
/usr/bin/pon
/etc/ppp/ppp_on_boot
/etc/ppp/ip-up
/etc/ppp/ip-up.d/
/etc/ppp/options
/etc/chatscripts/
/etc/ppp/chap-secrets
/usr/sbin/chat
/etc/ppp/peers/
/usr/bin/poff
/lib/systemd/system/ppp@.service
/etc/ppp/peers/provider
