# host_tools/sd_image/rootfs/var/lib/opkg/info/rpcbind.list
/usr/sbin/rpcbind
/etc/rpcbind.conf
/lib/systemd/system/rpcbind.socket
/usr/sbin/rpcinfo
/lib/systemd/system/rpcbind.service
