# host_tools/sd_image/rootfs/var/lib/opkg/info/timestamp-service.list
/lib/systemd/system/timestamp.service
/usr/bin/load-timestamp.sh
